# note/2022-07-13 - 守护进程/mysql.ini


[program:mysql]  # 这个是进程的名字，随意起
command=/usr/local/mysql/bin/mysqld --defaults-file=/etc/my.cnf
numprocs=1
stdout_logfile=/etc/supervisor/log/mysql.log
stderr_logfile=/etc/supervisor/log/mysqld.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 1200
priority=15
