# note/2022-07-13 - 守护进程/redis.ini
[program:redis]  # 这个是进程的名字，随意起
command=/usr/local/bin/redis-server  /home/redis/6379/conf/redis.conf
numprocs=1
stdout_logfile=/etc/supervisor/log/redis.log
stderr_logfile=/etc/supervisor/log/redisd.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 1200
priority=15
